# supervisord/vlm-service.conf
[program:vlm-service]
command=/root/gpufree-data/vlm-service/start.sh
directory=/root/gpufree-data/vlm-service
autostart=true
autorestart=true
startretries=20
startsecs=90
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=120
stdout_logfile=/root/gpufree-data/vlm-service/supervisor.out.log
stderr_logfile=/root/gpufree-data/vlm-service/supervisor.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
priority=20
